# idm-logger.service
[Unit]
Description=IDM Heat Pump Logger
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/idm-logger
ExecStart=/opt/idm-logger/venv/bin/python -m idm_logger.logger
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
